# systemd/rebootmgr.service
[Unit]
Description=Reboot Manager
Documentation=man:rebootmgrd(8) man:rebootmgrctl(1) man:rebootmgr.conf(5)
After=local-fs.target

[Service]
Type=Notify
ExecStart=/usr/libexec/rebootmgrd --verbose
Restart=on-failure

[Install]
WantedBy=multi-user.target

